# console.cron
*/3 * * * * curl --silent --show-error  http://localhost:3000/api/admin/events-log-trim
